# config/supervisor.conf
[program:rails_app]
directory=/app
command=/usr/local/bundle/bin/foreman start
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
